# .config/autostart/evolution.desktop
[Desktop Entry]
Name=Evolution
Exec=sh -c "/usr/bin/evolution"
Terminal=false
Type=Application
