# .config/autostart/firefox.desktop
[Desktop Entry]
Name=Firefox
Exec=sh -c "/usr/bin/firefox"
Terminal=false
Type=Application
