# .config/autostart/gedit.desktop
[Desktop Entry]
Name=Gedit
Exec=sh -c "/usr/bin/gedit"
Terminal=false
Type=Application
