# .config/autostart/keepassxc.desktop
[Desktop Entry]
Name=Keepassxc
Exec=sh -c "/usr/bin/keepassxc"
Terminal=false
Type=Application
