# .config/autostart/rclone.desktop
[Desktop Entry]
Name=Rclone
Exec=sh -c "rclone mount --daemon HetznerCR: /mnt/rclone"
Terminal=false
Type=Application
