# .config/autostart/standard_notes.desktop
[Desktop Entry]
Name=Standard Notes
Exec=sh -c "flatpak run org.standardnotes.standardnotes"
Terminal=false
Type=Application
