# .config/autostart/standard_notes.desktop
[Desktop Entry]
Name=Standard Notes
Exec=sh -c "/opt/standard-notes*"
Terminal=false
Type=Application
