# .config/autostart/syncthing.desktop
[Desktop Entry]
Name=Syncthing
Exec=sh -c "/usr/bin/syncthing"
Terminal=false
Type=Application
