# .config/autostart/terminal.desktop
[Desktop Entry]
Name=Terminal
Exec=sh -c "/usr/bin/gnome-terminal"
Terminal=false
Type=Application 
