# .config/autostart/vorta_flatpak.desktop
[Desktop Entry]
Name=Vorta Flatpak
Exec=sh -c "flatpak run com.borgbase.Vorta"
Terminal=false
Type=Application
